# SPAWN/STOP/.orchestrator/cron/parse_duplicates.cron
*/15 * * * * cd /mnt/d/NextAura/orchestrator && /usr/bin/python3 SPAWN/STOP/.orchestrator/maintenance.py --job queue_duplicate_parse_if_needed --runtime-dir /mnt/d/NextAura/orchestrator/SPAWN/STOP
